# ansible/files/servchal.service
[Service]
ExecStart=/opt/runChallenge.sh
Restart=always
User=root
Group=root
WorkingDirectory=/opt/
